# gitlab-runner/boot/com.chen.gitlab-runner-autostart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  GitLab Runner 開機自動啟動 LaunchAgent 範本。

  __RUNNER_DIR__ 為佔位符，由 install.sh 於安裝時以實際 gitlab-runner 目錄絕對路徑取代，
  並輸出至 ~/Library/LaunchAgents/com.chen.gitlab-runner-autostart.plist。
  請勿直接載入本範本檔。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chen.gitlab-runner-autostart</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__RUNNER_DIR__/boot/gitlab-runner-autostart.sh</string>
  </array>

  <!-- 使用者登入時即觸發一次（腳本內會等 Docker daemon 與 GitLab 就緒）。 -->
  <key>RunAtLoad</key>
  <true/>

  <!-- 標記為背景工作，避免佔用互動資源。 -->
  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>__RUNNER_DIR__/logs/gitlab-runner-autostart.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>__RUNNER_DIR__/logs/gitlab-runner-autostart.launchd.log</string>
</dict>
</plist>
